# deploy/supervisor/kanon-server.conf
[program:kanon-server]
command=/opt/kanon/scripts/run-kanon-server.sh
directory=/opt/kanon
autostart=true
autorestart=true
startsecs=2
stopsignal=TERM
stopasgroup=true
killasgroup=true
stdout_logfile=/var/log/kanon/server.stdout.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile=/var/log/kanon/server.stderr.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=3
environment=HOME="/home/kanon",USER="kanon",ADDR=":39090",ROOT="/root/docs",STATE_DIR="/var/lib/kanon/server",CONFIG_FILE="/opt/kanon/config/filter.json",RECONCILE_INTERVAL="1m",WATCH_DEBOUNCE="200ms",BUILD_ON_RUN="0"
