# deploy/supervisor/vault-bridge-server.conf
[program:vault-bridge-server]
command=/opt/vault-bridge/scripts/run-linux-server.sh
directory=/opt/vault-bridge
autostart=true
autorestart=true
startsecs=2
stopsignal=TERM
stopasgroup=true
killasgroup=true
stdout_logfile=/var/log/vault-bridge/server.stdout.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile=/var/log/vault-bridge/server.stderr.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=3
environment=HOME="/home/vault-bridge",USER="vault-bridge",ADDR=":39090",ROOT="/srv/vault-bridge/source",STATE_DIR="/var/lib/vault-bridge/server",CONFIG_FILE="/opt/vault-bridge/config/filter.json",RECONCILE_INTERVAL="30m",WATCH_DEBOUNCE="200ms",BUILD_ON_RUN="0"
