# deploy/systemd/user/kanon-server.service
[Unit]
Description=Kanon server
After=default.target

[Service]
Type=simple
Environment=HOME=%h
Environment=ADDR=:39090
Environment=ROOT=/root/docs
Environment=STATE_DIR=%h/.local/state/kanon/server
Environment=CONFIG_FILE=/opt/kanon/config/filter.json
Environment=RECONCILE_INTERVAL=1h
Environment=WATCH_DEBOUNCE=200ms
Environment=BUILD_ON_RUN=0
ExecStart=/opt/kanon/scripts/run-kanon-server.sh
Restart=always
RestartSec=2

[Install]
WantedBy=default.target
